# docker/supervisord/ELK.ini
[supervisord]
nodaemon = true

[program:elasticsearch]
command=/bin/bash -c "exec ${ES_HOME}/bin/elasticsearch -DFOREGROUND"

[program:logstash]
command=/zach/logstash/bin/logstash -f /logstash/logstash.conf

[program:kibana]
command=/zach/kibana/bin/kibana