# docker/supervisord/mango.ini
[supervisord]
nodaemon = true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mangod]
command=/zach/mongodb/bin/mongod -f /mongodb/mongodb.conf
autorestart = unexpected
autostart = true