# docker/supervisord/mysql.ini
[supervisord]
nodaemon = true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=./zach/mysql/bin/mysqld_safe --user=mysql
autorestart = unexpected
autostart = true