# docker/supervisord/redis.ini
[supervisord]
nodaemon = true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=/zach/src/redis-server /redis/redis.conf
autorestart = unexpected
autostart = true